# stowables/aconfmgr/.config/aconfmgr/files/etc/systemd/system/powertop.service
[Unit]
Description=Powertop tunings

[Service]
Type=oneshot
ExecStart=/usr/bin/powertop

[Install]
WantedBy=multi-user.target
